# user.service
[[Unit]
Description=Create User Automatically Service

[Service]
ExecStart=/home/uzmah/Documents/Create_User/TpUtilisateur/user.sh

