# user.service
[Unit]
Description=Create User Automatically Service

[Service]
ExecStart=/home/uzmah/src/TpUtilisateur/user.sh

